# var/spool/cron/crontabs/telegramd
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (crontab.telegramd installed on Tue Sep 11 23:26:49 2018)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
0 12 * * * /home/telegramd/noonphoto.sh
*/10 * * * * /home/pi/camera/camera.sh
